# php-fpm/82/supervisord.conf
[unix_http_server]
file=/var/run/supervisord.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true

[program:crond]
command=/usr/sbin/crond -f -l 8
autostart=true
autorestart=true

[program:consumer_sync_user]
command=php /var/www/hoozie/bank/spark amqp:consume-users user.sync.queue.peso_bank
autostart=true
autorestart=true
stderr_logfile=/var/log/consumer_sync_user.err.log
stdout_logfile=/var/log/consumer_sync_user.out.log
user=www-data

[program:consumer_sync_user_status]
command=php /var/www/hoozie/bank/spark amqp:consume-users user.status.queue.peso_bank
autostart=true
autorestart=true
stderr_logfile=/var/log/consumer_user_status.err.log
stdout_logfile=/var/log/consumer_user_status.out.log
user=www-data